# WebPyRobot/conf/webpyrobot_uwsgi_supervisord.conf
[program:webpyrobot]
command=/usr/local/lib/python3.6/bin/uwsgi --ini /var/www/MaturePyRobots/WebPyRobot/conf/webpyrobot.ini